# configs/launchd/kalshi_polygon_ws.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.kalshi.polygon-ws</string>

	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/env</string>
		<string>python3</string>
		<string>-m</string>
		<string>kalshi_alpha.exec.collectors.polygon_ws</string>
	</array>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PYTHONPATH</key>
		<string>src</string>
	</dict>

	<key>WorkingDirectory</key>
	<string>/Users/mateobodon/Documents/Programming/Projects/kalshi-sys</string>

	<key>RunAtLoad</key>
	<true/>

	<key>KeepAlive</key>
	<true/>

	<key>StandardOutPath</key>
	<string>/Users/mateobodon/Library/Logs/kalshi_polygon_ws.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/mateobodon/Library/Logs/kalshi_polygon_ws.err.log</string>
</dict>
</plist>
